# news_relay/crontab
# Run at 9 AM KST (Docker container uses UTC)
0 0 * * * echo `date` > /proc/1/fd/1 2>/proc/1/fd/2
0 0 * * * /usr/local/bin/python3 /news_relay/execute.py > /proc/1/fd/1 2>/proc/1/fd/2
